# deploy/systemd/aiq-testnet-user-stream.service
[Unit]
Description=AI Quant Binance Testnet user-data evidence observer
After=network-online.target aiq-testnet-secrets.service
Wants=network-online.target
Requires=aiq-testnet-secrets.service
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/quantify/ai-quant-system
UMask=0077
ExecStart=/root/quantify/ai-quant-system/.venv/bin/python -m ai_quant.services.testnet_user_stream --api-key-file /run/ai-quant-secrets/binance-testnet-api-key --api-secret-file /run/ai-quant-secrets/binance-testnet-api-secret --repository-root /root/quantify/ai-quant-system --evidence-file /var/lib/ai-quant/evidence/testnet/user-stream/current/events.jsonl --state-file /var/lib/ai-quant/evidence/testnet/user-stream/current/state.json --lock-file /run/ai-quant/testnet-user-stream.lock --keepalive-interval-seconds 1800 --rotate-no-later-than-seconds 84600
Restart=always
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
