# deploy/systemd/aiq-testnet-campaign.service
[Unit]
Description=AI Quant continuous Binance Testnet strategy campaign
After=network-online.target aiq-testnet-secrets.service
Wants=network-online.target
Requires=aiq-testnet-secrets.service
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=/root/quantify/ai-quant-system
UMask=0077
ExecStart=/root/quantify/ai-quant-system/.venv/bin/python -m ai_quant.services.testnet_campaign --api-key-file /run/ai-quant-secrets/binance-testnet-api-key --api-secret-file /run/ai-quant-secrets/binance-testnet-api-secret --repository-root /root/quantify/ai-quant-system --telegram-token-file /root/aiq-user-inputs/notifications/secrets/telegram_bot_token --telegram-chat-ids-file /root/aiq-user-inputs/notifications/telegram_chat_ids --evidence-directory /var/lib/ai-quant/evidence/testnet/campaign/current --state-file /var/lib/ai-quant/evidence/testnet/campaign/current/state.json --lock-file /run/ai-quant/testnet-campaign.lock --symbols BTCUSDT,ETHUSDT,BNBUSDT,SOLUSDT,XRPUSDT --duration-seconds 0 --evaluation-interval-seconds 60 --trade-cooldown-seconds 0 --maximum-trades-per-day 100 --daily-net-loss-limit 1.00 --margin-budget 1 --maximum-net-loss-per-trade 1.00 --maximum-parallel-positions 5 --maximum-candidates-per-round 5 --signal-confirmation-rounds 3 --impulse-confirmation-rounds 1 --minimum-signal-quality-score 2.00 --minimum-directional-forecast-bps 2.00 --impulse-minimum-directional-forecast-bps 0.10 --continuation-minimum-directional-forecast-bps 2.00 --structure-substitute-minimum-directional-forecast-bps 3.00 --structure-substitute-minimum-trade-imbalance 0.75 --structure-substitute-minimum-secondary-flow 0.10 --minimum-target-feasibility-rate-15m 0.20 --impulse-minimum-target-feasibility-rate-15m 0.02 --minimum-estimated-net-target 0.10 --minimum-net-reward-risk-ratio 0.50 --impulse-minimum-net-reward-risk-ratio 0.15 --risk-sizing-slippage-bps 12.00 --maximum-entry-spread-bps 5.00 --minimum-trade-imbalance 0.25 --minimum-book-imbalance 0.03 --minimum-microprice-bps 0.10 --maximum-opposing-book-imbalance 0.05 --maximum-opposing-microprice-bps 0.25 --aggressive-notional-lookback-rounds 6 --minimum-aggressive-notional-samples 3 --minimum-aggressive-notional-ratio 2.00 --no-activity-filter-enabled --no-impulse-activity-filter-enabled --impulse-minimum-activity-ratio 0.10 --impulse-maximum-activity-ratio 10.00 --impulse-lookback-rounds 4 --impulse-minimum-momentum-bps 2.00 --impulse-maximum-momentum-bps 8.00 --impulse-minimum-breadth-count 3 --sustained-lookback-rounds 5 --sustained-minimum-momentum-bps 5.00 --pullback-minimum-bps 3.00 --pullback-resumption-bps 0.50 --pullback-maximum-bps 40.00 --pullback-setup-maximum-rounds 10 --signal-evidence-window-rounds 2 --continuation-minimum-breadth-count 4 --continuation-confirmation-rounds 1 --continuation-minimum-momentum-bps 4.00 --continuation-maximum-momentum-bps 15.00 --position-failed-followthrough-peak-bps 6.00 --position-adverse-invalidation-bps 10.00 --position-profit-protection-peak-bps 20.00 --position-profit-giveback-bps 10.00 --position-opposition-confirmation-rounds 1
Restart=always
RestartSec=15
TimeoutStopSec=90
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
